# apps/chau7-macos/LaunchAgent/com.chau7.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chau7</string>

  <key>ProgramArguments</key>
  <array>
    <!-- EDIT THIS PATH to where you place the built binary -->
    <string>/bin/sh</string>
    <string>-lc</string>
    <string>$HOME/Applications/Chau7/Chau7</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- Optional: override log path -->
    <!-- <key>AI_EVENTS_LOG</key><string>$HOME/.ai-events.log</string> -->
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>
</dict>
</plist>
